# deploy/yagoda-backend.user.service
[Unit]
Description=Yagoda backend (uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/yagoda_staff/yagoda-bot
EnvironmentFile=/home/yagoda_staff/yagoda-bot/.env
ExecStart=/home/yagoda_staff/yagoda-bot/.venv/bin/uvicorn src.web.app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
